# deploy/funnel-worker.service
[Unit]
Description=Funnel Intel Worker
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/funnel-intel
EnvironmentFile=/opt/funnel-intel/.env
ExecStart=/opt/funnel-intel/.venv/bin/python -m backend.worker.loop
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
