# deploy/funnel-worker@.service
[Unit]
Description=Funnel Intel Worker (instance %i)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/funnel-intel
EnvironmentFile=/opt/funnel-intel/.env
Environment=WORKER_ID=%i
ExecStart=/opt/funnel-intel/.venv/bin/python -m backend.worker.loop
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
